# etc/supervisord.conf
[supervisord]
nodaemon=false
user=broker

[program:vnc]
environment=HOME="/home/broker",DISPLAY=":0",USER="broker"
command=tigervncserver --I-KNOW-THIS-IS-INSECURE -localhost no -fg :0 -SecurityTypes None,TLSNone
autorestart=true
priority=200

[program:tint2]
user=broker
environment=DISPLAY=":0"
command=tint2 -c /etc/xdg/tint2/tint2rc
autorestart=true
priority=300

[program:websockify]
user=broker
command=websockify -D --web /usr/share/novnc 6080 localhost:5900
autorestart=true
priority=400


[program:exec]
user=broker
command=bash -c "export EXEC_MODE=%(ENV_EXEC_MODE)s  && /home/broker/exec.sh"
priority=500
